# microauth.ini
[program:microauth]
command=gunicorn microauth.wsgi:application --bind localhost:8080
directory=/var/microauth
autostart=true
autorestart=true
startretries=3
