# microauth.ini
[program:microauth]
command=gunicorn microauth.wsgi:application --bind localhost:8080 --log-file /var/log/microauth_guni.log
stdout_logfile=/var/log/stdout-microauth.log
directory=/var/microauth
autostart=true
autorestart=true
startretries=3
